# templates/puppet-run.timer.erb
[Unit]
Description=Schedule puppet run

[Timer]
OnBootSec=<%= @bootdelay %>
OnUnitActiveSec=<%= @frequency %>

[Install]
WantedBy=multi-user.target

